# bin/default/packaging/systemd/nmb.service
[Unit]
Description=Samba NMB Daemon
Documentation=man:nmbd(8) man:samba(7) man:smb.conf(5)
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=notify
NotifyAccess=all
PIDFile=/Users/loctv/Documents/QEMU/raspian/samba/build/var/run/nmbd.pid
EnvironmentFile=-/Users/loctv/Documents/QEMU/raspian/samba/build/etc/sysconfig/samba
ExecStart=/Users/loctv/Documents/QEMU/raspian/samba/build/sbin/nmbd --foreground --no-process-group $NMBDOPTIONS
ExecReload=/bin/kill -HUP $MAINPID
LimitCORE=infinity


[Install]
WantedBy=multi-user.target
