# bin/default/packaging/systemd/smb.service
[Unit]
Description=Samba SMB Daemon
Documentation=man:smbd(8) man:samba(7) man:smb.conf(5)
Wants=network-online.target
After=network.target network-online.target nmb.service winbind.service

[Service]
Type=notify
NotifyAccess=all
PIDFile=/Users/loctv/Documents/QEMU/raspian/samba/build/var/run/smbd.pid
LimitNOFILE=16384
EnvironmentFile=-/Users/loctv/Documents/QEMU/raspian/samba/build/etc/sysconfig/samba
ExecStart=/Users/loctv/Documents/QEMU/raspian/samba/build/sbin/smbd --foreground --no-process-group $SMBDOPTIONS
ExecReload=/bin/kill -HUP $MAINPID
LimitCORE=infinity


[Install]
WantedBy=multi-user.target
