# bin/default/packaging/systemd/winbind.service
[Unit]
Description=Samba Winbind Daemon
Documentation=man:winbindd(8) man:samba(7) man:smb.conf(5)
After=network.target nmb.service

[Service]
Type=notify
NotifyAccess=all
PIDFile=/opt/samba/var/run/winbindd.pid
EnvironmentFile=-/opt/samba/etc/sysconfig/samba
ExecStart=/opt/samba/sbin/winbindd --foreground --no-process-group "$WINBINDOPTIONS"
ExecReload=/bin/kill -HUP $MAINPID
LimitCORE=infinity


[Install]
WantedBy=multi-user.target
